# 99-finepix.rules
ATTRS{idVendor}=="04cb", ATTRS{idProduct}=="0186", SUBSYSTEMS=="usb" KERNEL=="sda1", RUN+="/home/stamp/finepix.sh"

